# parser-service/crontab
0 3 * * * cd /app && /usr/local/bin/python parser.py >> /var/log/parser.log 2>&1
